# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/usr/src/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:hopper]
command=/usr/src/app/docker/hopper-gunicorn.sh
directory=/usr/src/app
autostart=true
autorestart=true
stdout_logfile=/usr/src/app/logs/hopper_supervisor.log
stderr_logfile=/usr/src/app/logs/hopper_supervisor.log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true
stopasgroup=true

priority=998
